# config/systemd/highlightly-match-lifecycle-report.service
[Unit]
Description=ASP Insights Highlightly Phase 8E.1 operational report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/ubuntu/asp-insights-c73dbc6b
EnvironmentFile=/etc/asp-scraper-api.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/ubuntu/asp-scraper-api/.venv/bin/python -m scripts.report_highlightly_phase8e_operational --hours 24 --require-provider-disabled
TimeoutStartSec=120
Nice=10
