# config/systemd/highlightly-match-lifecycle.service
[Unit]
Description=ASP Insights Highlightly Phase 8E match lifecycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/ubuntu/asp-insights-c73dbc6b
EnvironmentFile=/etc/asp-scraper-api.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/flock --exclusive --nonblock --conflict-exit-code 0 /run/lock/asp-highlightly-future.lock /home/ubuntu/asp-scraper-api/.venv/bin/python -m scripts.run_highlightly_phase8e_match_lifecycle --request-budget 1500 --max-jobs 1000 --confirm-lifecycle
TimeoutStartSec=infinity
Nice=0
