# config/systemd/highlightly-future-continuation.service
[Unit]
Description=ASP Insights Highlightly future queue continuation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/ubuntu/asp-insights-c73dbc6b
EnvironmentFile=/etc/asp-scraper-api.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/flock --exclusive --nonblock --conflict-exit-code 0 /run/lock/asp-highlightly-future.lock /home/ubuntu/asp-scraper-api/.venv/bin/python -m scripts.run_highlightly_future_continuation --max-jobs 5000 --confirm-continuation
TimeoutStartSec=infinity
Nice=5
